# docker/supervisor.d/host1.ini
[program:normal1]
command=/app/normal.py

[program:normal2]
command=/app/normal.py
process_name=%(program_name)s_%(process_num)02d
numprocs=2

[group:group_normal]
programs=normal1,normal2

[program:crash_now]
command=/app/crash_now.py

[program:crash_15s]
command=/app/crash_15s.py

[group:group_crash]
programs=crash_now,crash_15s

[program:exit_2]
command=/app/exit_2.py

[program:run_1_minute]
command=/app/run_1_minute.py

[program:ignore_signal]
command=/app/ignore_signal.py
