# docker/supervisor.d/host3.ini
[program:normal1]
command=/app/normal.py

[program:normal2]
command=/app/normal.py

[group:group_normal]
programs=normal1,normal2

[program:exit_2]
command=/app/exit_2.py

[program:run_1_minute]
command=/app/run_1_minute.py

[program:ignore_signal]
command=/app/ignore_signal.py
